# certbot-renew.timer
[Unit]
Description=Weekly Certbot renewal

[Timer]
OnCalendar=Weekly
Persistent=true

[Install]
WantedBy=multi-user.target
